# deploy/mqtt-live-monitor.service
[Unit]
Description=MQTT Live Monitor
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/mqtt-live-monitor
ExecStart=/usr/bin/npm run start
Restart=always
RestartSec=5
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target